# root/etc/services.d/hagrid/run
#!/bin/bash

cd /data

exec s6-setuidgid abc /usr/local/bin/hagrid
